# admin_scripts/security/91-detectusbevents.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/etc/os2borgerpc/security/keyboardEvent.py 'input:'"
